# sources/meta-multiboot-update/recipes-platform/update-observer/files/systemd/update-observer.service
[Unit]
Description=Initiates the software update.

[Service]
Type=simple
ExecStart=/usr/sbin/update-observer

[Install]
WantedBy=multi-user.target