# layers/meta-multiboot-update/recipes-platform/update-observer/files/systemd/regular-operation.target
[Unit]
Description=Indicates regular operation
Wants=update-observer.service
After=update-observer.service

[Install]
WantedBy=default.target

